# openfire/supervisord/openfire.sv.conf
[program:openfire]
command=/bin/bash -c "service openfire start"
redirect_stderr=true
